# deploy/server/systemd/quantpilot-web.service
[Unit]
Description=QuantPilot Next.js web application
Wants=network-online.target
After=network-online.target docker.service quantpilot-market-data.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/quantpilot
EnvironmentFile=/opt/quantpilot/.env.production
Environment=HOME=/home/ubuntu
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin:/snap/bin
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGINT
UMask=0027

[Install]
WantedBy=multi-user.target
